# misc/95-unipi-plc-devices.rules
########################################################################
# Copyright (C) 2009-2011 Richard Hughes <richard@hughsie.com>
#
# Licensed under the GNU General Public License Version 2
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation; either version 2 of the License, or
# (at your option) any later version.
#
########################################################################

# SPI UART
KERNEL=="ttyNS0", SYMLINK+="extcomm/0/0"
KERNEL=="ttyNS0", SYMLINK+="extcomm/1/0", KERNELS=="spi0.1"
KERNEL=="ttyNS0", SYMLINK+="extcomm/2/0", KERNELS=="spi0.2"
KERNEL=="ttyNS0", SYMLINK+="extcomm/3/0", KERNELS=="spi0.3"

KERNEL=="ttyNS1", SYMLINK+="extcomm/0/1"
KERNEL=="ttyNS1", SYMLINK+="extcomm/1/0", KERNELS=="spi0.1"
KERNEL=="ttyNS1", SYMLINK+="extcomm/2/0", KERNELS=="spi0.2"
KERNEL=="ttyNS1", SYMLINK+="extcomm/3/0", KERNELS=="spi0.3"

KERNEL=="ttyNS2", SYMLINK+="extcomm/0/2"
KERNEL=="ttyNS2", SYMLINK+="extcomm/1/0", KERNELS=="spi0.1"
KERNEL=="ttyNS2", SYMLINK+="extcomm/2/0", KERNELS=="spi0.2"
KERNEL=="ttyNS2", SYMLINK+="extcomm/3/0", KERNELS=="spi0.3"

# SYSFS Platform Device Symlink
KERNEL=="unipi_plc", RUN+="/bin/sh -c 'mkdir -p /run; ln -s /sys/devices/platform/unipi_plc /run/unipi 2>/dev/null|| exit 0'"
